# salt/files/jenkins/supervisor/dashing.conf
[program:dashing]
command=dashing start
autostart=true
directory=/srv/dashboard
